# systemd/skypier-blackhole.service
[Unit]
Description=Skypier Blackhole DNS Resolver
Documentation=https://github.com/skypier/skypier-blackhole
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=skypier
Group=skypier

# Main process
ExecStart=/usr/bin/skypier-blackhole start --config /etc/skypier/blackhole.toml

# Reload blocklists without restart (SIGHUP)
ExecReload=/bin/kill -HUP $MAINPID

# Restart policy
Restart=on-failure
RestartSec=5s

# Timeout for start/stop
TimeoutStartSec=30s
TimeoutStopSec=10s

# Resource limits
LimitNOFILE=65536
LimitNPROC=512

# Working directory
WorkingDirectory=/var/lib/skypier

# Environment
Environment="RUST_LOG=info"

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/skypier /etc/skypier
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=true
LockPersonality=true
RestrictRealtime=true
RestrictSUIDSGID=true
RemoveIPC=true
PrivateMounts=true

# Network capabilities (required for binding to port 53)
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

# Systemd features
StandardOutput=journal
StandardError=journal
SyslogIdentifier=skypier-blackhole

[Install]
WantedBy=multi-user.target
